# debian/lte/lte-connect.service
[Unit]
Description=LTE Auto Connect via ModemManager
After=ModemManager.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/lte-connect.sh
Restart=on-failure
RestartSec=15
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target